# run/triangle.conf
[program:triangle]
directory=/home/dev/vulkan
command=/home/dev/vulkan/triangle
redirect_stderr=true
