# run/xfce.conf
[program:xfce]
command=xfce4-session
redirect_stderr=true
